# deploy/spikelab-daemon.service
# OC-SpikeLab MCP daemon — systemd USER unit.
#
# Default install path (customer device): runs as the user that owns this
# install. The customer's existing $HOME, $XDG_RUNTIME_DIR, and
# ~/.claude/.credentials.json are used as-is — no dedicated UNIX account
# required.
#
# Install (customer-device, no sudo for the unit itself):
#
#   mkdir -p ~/.config/systemd/user
#   cp <repo>/deploy/spikelab-daemon.service ~/.config/systemd/user/
#
#   # Point ExecStart at whichever Python interpreter has the [mcp] extra
#   # installed (analysis env, e.g. the OC_env conda env). The default
#   # below uses /usr/bin/env python which only works if the [mcp] extra
#   # is on the default PATH for this systemd-user environment.
#   sed -i 's|^ExecStart=.*|ExecStart=/home/<user>/anaconda3/envs/OC_env/bin/python -m spikelab.mcp.daemon|' \
#       ~/.config/systemd/user/spikelab-daemon.service
#
#   # Tell the daemon where the project root (system_config.json) lives.
#   # systemd --user automatically reads environment.d:
#   mkdir -p ~/.config/environment.d
#   echo "OC_SPIKELAB_PROJECT_ROOT=/path/to/Habitat-MaxOne-test" \
#       > ~/.config/environment.d/oc-spikelab.conf
#
#   systemctl --user daemon-reload
#   systemctl --user enable --now spikelab-daemon.service
#
#   # ONE sudo step: keep the unit alive when the user isn't logged in
#   # interactively. Required if you want the daemon to survive logouts.
#   sudo loginctl enable-linger $USER

[Unit]
Description=OC-SpikeLab MCP daemon (persistent per-task_id child agents)
After=default.target

[Service]
Type=simple
# Default: rely on PATH (works if the [mcp]-extra Python is on PATH).
# Most installs will override this to point at their analysis-env python,
# e.g. /home/<user>/anaconda3/envs/OC_env/bin/python — see install docs.
ExecStart=/usr/bin/env python -m spikelab.mcp.daemon

# OC_SPIKELAB_PROJECT_ROOT is NOT set here. Provide it via:
#   - ~/.config/environment.d/oc-spikelab.conf (preferred, systemd-user reads
#     this automatically), or
#   - A drop-in override that adds `Environment=OC_SPIKELAB_PROJECT_ROOT=...`

# Keep the daemon alive across crashes; back off briefly to avoid restart
# storms if something is fundamentally broken.
Restart=always
RestartSec=2s

StandardOutput=journal
StandardError=journal

# Mild hardening. The daemon needs write access to its socket path and to
# /tmp; otherwise nothing privileged.
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=default.target
